# shared/pairing_sync.service.template
[Unit]
Description=Pairing Synchronization Service
After=pairing_deployment.target

[Service]
Type=simple
ExecStart=${DIR}/sync.sh
User=${SCRIPT_USER}
Restart=always
RestartSec=60
